# mouseless.service
[Unit]
Description=mouseless service

[Service]
ExecStart=/usr/bin/mouseless

[Install]
WantedBy=default.target
