# mouseless.service
[Unit]
Description=mouseless

[Service]
ExecStartPre=/bin/sleep 2
ExecStart=/usr/local/bin/mouseless

[Install]
WantedBy=default.target
